# dot_config/systemd/user/gnome-polkit-agent.service
[Unit]
Description=GNOME polkit authentication agent
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=graphical-session.target
